# sensu-server/sensu/config/supervisord.conf
[supervisord]
nodaemon=true

[program:sensu-server]
command=/etc/init.d/sensu-server start

[program:uchiwa]
command=/etc/init.d/uchiwa start

[program:sensu-api]
command=/etc/init.d/sensu-api start
